# scripts/server/crontab
MAILTO=torsten@metaluna4.de
*/5 * * * *		~/cubestats/scripts/server/cron_daemoncontroller.sh

